# meta-video-surveillance/recipes-core/usb/files/99-usb-automount.rules
ACTION=="add", KERNEL=="sd*[0-9]", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/lib/udev/scripts/usb-automount.sh %k"
